# cron.d/crontab-default
# Here's an example crontab you can adapt for /etc/cron.d on your system.

# Startup error messages are directed to stderr, so it's safe to dump stdout here.
# Change "ictmanager" to the owner of your target volumes.
0 18 * * * ictmanager /srv/lime-machine/scripts/start-backup.sh > /dev/null

